# deploy/penumbra.service
[Unit]
Description=PENUMBRA Space-Weather Driver Forecasting
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/penumbra/backend
ExecStart=/opt/penumbra/backend/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8001
Restart=on-failure
RestartSec=10
MemoryMax=190M
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
